# deploy/r3p-backend.service
[Unit]
Description=R3P FastAPI backend
After=network-online.target r3p-db.service
Requires=r3p-db.service

[Service]
Type=simple
Restart=always
RestartSec=10
ExecStart=/usr/bin/podman run --rm --name r3p-backend \
  --env-file /home/admin/r3p/.env \
  -v /home/admin/r3p/logs:/app/logs:Z \
  -p 127.0.0.1:8000:8000 \
  r3p-backend:latest
ExecStop=/usr/bin/podman stop -t 10 r3p-backend

[Install]
WantedBy=default.target
